# nextcloud.cron
# /etc/cron.d/nextcloud: crontab entries for nextcloud
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/5 * * * * www-data php -f /var/www/nextcloud/cron.php
30 5 * * * root /bin/05_run_app_update.sh
